# panel/systemd/tenodera-gateway.service
[Unit]
Description=Tenodera Web Console Gateway
Documentation=https://github.com/tenodera
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tenodera-gw
Group=tenodera-gw
ExecStart=/usr/local/bin/tenodera-gateway
EnvironmentFile=-/etc/tenodera/gateway.env

# Security hardening
#
# The gateway runs as the unprivileged tenodera-gw user.
# tenodera-pam-helper and tenodera-bridge are setuid root (4750, root:tenodera-gw)
# so they can escalate for PAM auth and PTY user switching respectively.
#
# NoNewPrivileges is intentionally NOT set: it would prevent the setuid
# pam-helper and bridge from gaining root, breaking authentication and PTY.
#
# ProtectSystem=strict is NOT used: child processes (bridge) write to /usr,
# /var, /run for package management, systemd, docker etc.
#
# MemoryDenyWriteExecute is NOT used: PAM modules (pam_sss.so) may require
# W+X mappings.
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
RestrictNamespaces=yes
LockPersonality=yes
# pam-helper and bridge are intentionally setuid root
RestrictSUIDSGID=no
# Gateway itself needs no special capabilities — privileged ops go through
# the setuid helper binaries. Bounding set prevents the gateway process from
# acquiring any capabilities even if exploited.
CapabilityBoundingSet=
AmbientCapabilities=

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
