# panel/systemd/tenodera-gateway.service
[Unit]
Description=Tenodera Web Console Gateway
Documentation=https://github.com/tenodera
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tenodera-gateway
EnvironmentFile=-/etc/tenodera/gateway.env

# The gateway starts as root so it can read TLS cert/key files from their
# standard system location (e.g. /etc/ssl/...) regardless of permissions,
# then immediately drops privileges to the tenodera-gw user via setuid(2)
# after the socket is bound and certs are loaded — same pattern as nginx,
# Zabbix, and Apache.
#
# NoNewPrivileges is intentionally NOT set: child processes (tenodera-pam-helper,
# tenodera-bridge) are setuid root and must be able to gain root for PAM auth
# and PTY user switching.
#
# The following systemd security directives are intentionally NOT used
# because they break setuid child processes (pam-helper, bridge):
#
#   PrivateTmp=yes            — user namespace mount isolation drops setuid
#   RestrictNamespaces=yes    — implies NoNewPrivileges
#   LockPersonality=yes       — implies NoNewPrivileges
#   ProtectKernelTunables=yes — implies NoNewPrivileges on this systemd version
#   ProtectKernelModules=yes  — implies NoNewPrivileges on this systemd version
#   CapabilityBoundingSet=    — empty set inherited by children; breaks setuid helpers
#   MemoryDenyWriteExecute    — PAM modules (pam_sss.so, pam_ldap.so) may need W+X
#   ProtectSystem=strict      — bridge child processes write to /usr, /var, /run
RestrictSUIDSGID=no
ProtectControlGroups=yes

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
